# config/crontab/backup.sh
#!/bin/bash
mysqldump --user=farmacentro --password=farmacentro Farmacentro > /home/hditx/fmaci/config/data.sql

