# config/crontab/backup.sh
#!/bin/bash
mysqldump --user=farmacentro --password=farmacentro Farmacentro > /media/firefly/fire1/backup/data.sql

