# config/crontab/cronfarmacia
50 * * * * sh /home/firefly/fmaci/config/crontab/backup.sh
7 * * * * php /home/firefly/fmaci/config/crontab/base.php

