# config/crontab/cronfarmacia
50 * * * * sh /home/hditx/prueba.sh
7 * * * * php /home/hditx/base.php

