# config/crontab/cronfarmacia
30 23 * * * sh /media/firefly/fire1/fmaci/config/crontab/backup.sh
30 23 * * * php /media/firefly/fire1/fmaci/config/crontab/base.php
1 * * * * touch crear.txt

